# src/etc/services.d/tracks/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid tracks /opt/bin/run.sh
