# src/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid tracks /usr/sbin/nginx -c /etc/nginx/nginx.conf -g 'daemon off;' -p /data/nginx
